# support/backup.crontab
# Helpful reference: http://www.thegeekstuff.com/2009/06/15-practical-crontab-examples
#
# "every day at 1 AM run the backup script"
# minutes hour dayOfMonth month dayOfWeek command
  00      01   *          *     *         /bin/bash -c '/export/support/run_backup.sh >> /export/var/log/run_backup.log 2>&1'
